# ops/launchd/com.giuseppelopesme.openclaw.brain.clu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  OpenClaw CLU brain launchd job.

  Runs as macOS user `giuseppelopes` (NOT `clu` — only the iMessage
  relay runs as `clu`; brains run as the operator). The brain
  subscribes to `imessage.received.clu` over the bridge's WebSocket
  (loopback), runs LLM triage + draft, and publishes
  `agent.clu.draft.pending` for the human-approval flow that lands
  in a future session.

  Pre-flight (operator, one-time):

  1. Mint a brain token (stores into giuseppelopes's Keychain; the
     wrapper script reads it back at launch — no plaintext in this
     plist, no env var to substitute):
       scripts/mint-token.py \
           --actor brain.clu \
           --scopes llm:call,vault:read,vault:write,events:subscribe,events:publish,agent:drafts:write

  2. Make sure the bridge + Redis are running first. With the bridge,
     redis, brain.clu, and relay.clu plists all installed, launchd
     starts everything on boot.

  NOT auto-installed by Session 9. To load (one-time):

      cp ops/launchd/com.giuseppelopesme.openclaw.brain.clu.plist \
         ~/Library/LaunchAgents/
      launchctl load ~/Library/LaunchAgents/com.giuseppelopesme.openclaw.brain.clu.plist

  To unload:
      launchctl unload ~/Library/LaunchAgents/com.giuseppelopesme.openclaw.brain.clu.plist

  Logs are written to /Users/giuseppelopes/.openclaw/clu.{out,err}.log
  so Console.app and `tail -f` both work. The brain token is loaded
  from Keychain by scripts/run-clu.sh — no secret in this plist.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.giuseppelopesme.openclaw.brain.clu</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/giuseppelopes/Runtime/OpenClaw_Bridge/scripts/run-clu.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/giuseppelopes/Runtime/OpenClaw_Bridge</string>

  <key>StandardOutPath</key>
  <string>/Users/giuseppelopes/.openclaw/clu.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/giuseppelopes/.openclaw/clu.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>

    <key>BRIDGE_URL</key>
    <string>http://127.0.0.1:8788</string>

    <key>STATE_DB_PATH</key>
    <string>/Users/giuseppelopes/.openclaw/clu.state.db</string>
  </dict>
</dict>
</plist>
