# ops/launchd/com.giuseppelopesme.openclaw.relay.clu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  OpenClaw iMessage relay (CLU) launchd job.

  Runs as macOS user `clu`, NOT `giuseppelopes`. CLU has its own
  ~/Library/Messages/chat.db (per-user) and its own iMessage account
  signed-in via Messages.app.

  Pre-flight (operator, one-time):

  1. macOS user `clu` exists. Create via System Settings → Users & Groups
     → Add User if not.
  2. Sign into iMessage as that user (Messages.app).
  3. Grant Full Disk Access to /usr/bin/python3 (or the relay's actual
     binary) so it can read chat.db. System Settings → Privacy & Security
     → Full Disk Access → Add → /usr/bin/python3.
  4. Grant Automation control of Messages.app on first send. The
     "Terminal/Python wants to control Messages.app" prompt fires the
     first time; click OK.
  5. Mint a relay token on the bridge host (giuseppelopes), then store
     the same plaintext in clu's login keychain so scripts/run-relay.sh
     can read it. Two ways:
       a) On giuseppelopes, mint and capture the plaintext:
            scripts/mint-token.py --actor relay.clu --scopes imessage:relay
          Then on clu, open Keychain Access, add a generic password with
          service=`com.giuseppelopesme.openclaw.bridge`, account=`relay.clu`,
          password=`{"token":"<plaintext>","scopes":["imessage:relay"]}`.
       b) Or just rerun mint-token.py from the clu account (it stores
          into the local keychain). Bridge token validation happens on
          giuseppelopes's side, so you must ALSO have the entry in
          giuseppelopes's keychain — option (a) is normally simpler.

  NOT auto-installed by Session 7. To load (one-time, from the `clu`
  account — note: `clu` needs its own production runtime clone at
  /Users/clu/Runtime/OpenClaw_Bridge because launchd can't read into
  another user's home dir):

      sudo -u clu git clone https://github.com/giuseppelopesme/OpenClaw_Bridge.git \
         /Users/clu/Runtime/OpenClaw_Bridge
      sudo -u clu sh -c 'cd /Users/clu/Runtime/OpenClaw_Bridge && uv sync --group dev'
      sudo cp ops/launchd/com.giuseppelopesme.openclaw.relay.clu.plist \
         /Users/clu/Library/LaunchAgents/
      sudo chown clu:staff /Users/clu/Library/LaunchAgents/com.giuseppelopesme.openclaw.relay.clu.plist
      sudo -u clu launchctl bootstrap gui/$(id -u clu) \
         /Users/clu/Library/LaunchAgents/com.giuseppelopesme.openclaw.relay.clu.plist

  To unload:
      launchctl unload ~/Library/LaunchAgents/com.giuseppelopesme.openclaw.relay.clu.plist

  Logs are written to /Users/clu/.openclaw/relay.clu.{out,err}.log so
  Console.app and `tail -f` both work. The relay token is loaded from
  clu's login keychain by scripts/run-relay.sh — no secret in this plist.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.giuseppelopesme.openclaw.relay.clu</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/clu/Runtime/OpenClaw_Bridge/scripts/run-relay.sh</string>
    <string>clu</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/clu/Runtime/OpenClaw_Bridge</string>

  <key>StandardOutPath</key>
  <string>/Users/clu/.openclaw/relay.clu.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/clu/.openclaw/relay.clu.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>

    <key>AGENT_NAME</key>
    <string>clu</string>

    <key>BRIDGE_URL</key>
    <string>http://127.0.0.1:8788</string>
  </dict>
</dict>
</plist>
